# service.bookmarking/roles/service.bookmarking.readeck.backup/templates/readeck-backup.service
[Unit]
Description=restic backup of readeck data

[Service]
Type=oneshot
ExecStart={{ ansible_env.HOME }}/.config/restic/readeck/readeck-backup.sh
