# notesync-unison/ansible/roles/peer.unison.sync/templates/unison-sync.timer.j2
[Unit]
Description=Synchronize files with unison

[Timer]
OnCalendar=*-*-* 8,15:00:00
Persistent=true

[Install]
WantedBy=timers.target
