# ansible/roles/server.bookmarking.readeck.backup.replica/templates/readeck-backup-replicate.service
[Unit]
Description=Replicate readeck restic snapshots

[Service]
Type=oneshot
ExecStart={{ readeck_restic_config_dir }}/readeck-backup-replicate.sh
