# ansible/roles/server.bookmarking.readeck.backup.replica/files/readeck-backup-replicate.timer
[Unit]
Description=Replicate readeck restic snapshots

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
