# ansible/roles/server.srv.backup/templates/srv-backup.service
[Unit]
Description=restic backup of srv

[Service]
Type=oneshot
ExecStart={{ ansible_env.HOME }}/.config/restic/srv/srv-backup.sh
