# ansible/roles/hub.restic.replica/files/restic-disk-replica.timer
[Unit]
Description=Copy restic snapshots trigger

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
