# ansible/roles/server.zettelkasten.git.backup.replica/templates/zettelkasten-git-replica.service
[Unit]
Description=Copy restic snapshots to replica

[Service]
Type=oneshot
ExecStart={{ zettelkasten_restic_config_dir }}/zettelkasten-git-replica.sh
