# notesync-git/ansible/roles/peer.git.sync/files/zettelkasten-git-sync.timer
[Unit]
Description=Zettelkasten git sync timer

[Timer]
# daily at 00:00:00, according systemd.time(7)
OnCalendar=daily
# run ASAP if last schedule did not run
Persistent=true
# Run 5min after boot
OnBootSec=300

[Install]
WantedBy=timers.target
